# deploy/roles/webapp/templates/systemd.webapp.socket
[Unit]
Description=gunicorn webapp socket

[Socket]
ListenStream={{ webapp_socket_path }}

[Install]
WantedBy=sockets.target
